# bin/itc-monitor
#!/bin/sh

SCRIPT=$(realpath "$0")
DIR=$(dirname "$SCRIPT")
CMD="${DIR}/itc-monitor.sh"
NAME="itc-rest"
RUNAS=root
DESCRIPTION="ITC Monitor Agent Service"

### BEGIN INIT INFO
# Provides:         $NAME
# Required-Start:	$local_fs $network $named $time $syslog
# Required-Stop:	$local_fs $network $named $time $syslog
# Default-Start:	2 3 4 5
# Default-Stop:		0 1 6
# Description:      $DESCRIPTION
### END INIT INFO

source "$DIR/itc-initd" $1